# deploy/jat-bot.service.example
[Unit]
Description=JAT Telegram Assistant Bot
After=network.target

[Service]
Type=simple
User=<LINUX_USER>
WorkingDirectory=/opt/jat_telegram_assistant_bot
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/jat_telegram_assistant_bot/.venv/bin/uvicorn app:app --host 0.0.0.0 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

